# debian/minios-live-config-systemd.postinst
#!/bin/sh

set -e

case "$1" in
configure)
    if command -v deb-systemd-helper >/dev/null; then
        deb-systemd-helper enable live-config.service >/dev/null || true
    fi
    ;;
esac

#DEBHELPER#

exit 0
